# config/RPi-Edge-A/etc/systemd/system/speedtest-web.service
# container-0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88.service
# autogenerated by Podman 3.4.1
# Wed Dec 15 19:45:06 EST 2021

[Unit]
Description=Podman container-0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88
ExecStop=/usr/bin/podman stop -t 10 0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88
ExecStopPost=/usr/bin/podman stop -t 10 0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88
PIDFile=/run/containers/storage/overlay-containers/0b4501fb2294352953acd2120319f858b50cee96ef3def813f7b6cca32075e88/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
