# config/RPi-Edge-A/etc/systemd/system/unifi.service
# container-18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f.service
# autogenerated by Podman 3.4.1
# Wed Dec 15 19:44:24 EST 2021

[Unit]
Description=Podman container-18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f
ExecStop=/usr/bin/podman stop -t 10 18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f
ExecStopPost=/usr/bin/podman stop -t 10 18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f
PIDFile=/run/containers/storage/overlay-containers/18873136a1090d80ab90bf253d1411427d47b9e6d2fea9cab6acdaa2bcdd983f/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
